# supervisord.conf
; supervisord config for Larynx workers.
; In M3 the gateway, VoxCPM2 worker, Fun-ASR worker, and VAD+Punc worker
; share a single Python process (via the in-process WorkerChannel) — so
; this file only declares the gateway program when everything runs on
; one box.
;
; When workers are split to dedicated containers (see
; docker/funasr_worker.Dockerfile + the `workers` compose profile), each
; container runs supervisord with only its respective [program:*] block.
; The `funasr_worker` + `vad_punc_worker` blocks below are left in one
; place so every container reads the same supervisord.conf with the
; right subset uncommented at image-build time.

[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid

[program:gateway]
command=uv run uvicorn larynx_gateway.main:app --host 0.0.0.0 --port 8000
directory=/app
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:funasr_worker]
command=uv run larynx-funasr-worker
directory=/app
autostart=false
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=LARYNX_STT_MODE=funasr,LARYNX_FUNASR_GPU=0,CUDA_VISIBLE_DEVICES=1

[program:vad_punc_worker]
command=uv run larynx-vad-punc-worker
directory=/app
autostart=false
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=LARYNX_VAD_PUNC_MODE=real
