# supervisord.conf
; supervisord config for Larynx workers.
; In M3 the gateway, VoxCPM2 worker, Fun-ASR worker, and VAD+Punc worker
; share a single Python process (via the in-process WorkerChannel) — so
; this file only declares the gateway program when everything runs on
; one box.
;
; When workers are split to dedicated containers (see
; docker/funasr_worker.Dockerfile + the `workers` compose profile), each
; container runs supervisord with only its respective [program:*] block.
; The `funasr_worker` + `vad_punc_worker` blocks below are left in one
; place so every container reads the same supervisord.conf with the
; right subset uncommented at image-build time.
;
; M8 §3.3 — restart policy:
; - ``startsecs=5`` means the process must stay alive 5s before
;   supervisord considers it "running"; crash-during-boot counts as
;   a restart attempt.
; - ``startretries=3`` gives us 3 rapid restarts before supervisord
;   moves the program to FATAL. Combined with the restart_alerter
;   event listener, 3 fails inside 60s emits a critical.restart_storm
;   line that the soak test + future pagerduty hook watch.
; - ``stopwaitsecs=35`` gives the lifespan 30s to drain in-flight
;   batch items + cleanup cron + WS connections before SIGKILL.

[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid

[program:gateway]
command=uv run uvicorn larynx_gateway.main:app --host 0.0.0.0 --port 8000
directory=/app
autostart=true
autorestart=true
startsecs=5
startretries=3
stopwaitsecs=35
stopsignal=TERM
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:funasr_worker]
command=uv run larynx-funasr-worker
directory=/app
autostart=false
autorestart=true
startsecs=5
startretries=3
stopwaitsecs=15
stopsignal=TERM
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=LARYNX_STT_MODE=funasr,LARYNX_FUNASR_GPU=0,CUDA_VISIBLE_DEVICES=1

[program:vad_punc_worker]
command=uv run larynx-vad-punc-worker
directory=/app
autostart=false
autorestart=true
startsecs=5
startretries=3
stopwaitsecs=15
stopsignal=TERM
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=LARYNX_VAD_PUNC_MODE=real

[eventlistener:restart_alerter]
; Watches PROCESS_STATE events; emits a critical.restart_storm line
; when a single program fails to start 3 times within 60s. Keeps
; supervisord's own log as the single source of alerting truth until
; we wire a pagerduty/prometheus alert.
command=uv run python -m larynx_gateway.ops.restart_alerter
events=PROCESS_STATE
autostart=true
autorestart=true
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
