# supervisord.conf
; supervisord config for Larynx workers (M1 placeholder).
; In M1 the gateway and VoxCPM2 worker share a single Python process, so this
; file only declares the gateway. Additional workers (funasr, vad_punc,
; training) land here in later milestones.

[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid

[program:gateway]
command=uv run uvicorn larynx_gateway.main:app --host 0.0.0.0 --port 8000
directory=/app
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
